# tt-leaderboard.service
[Unit]
Description=Elite Dangerous Time Trials Leaderboard
After=network.target

[Service]
Type=simple
User=elite
Group=elite
WorkingDirectory=/srv/elite/timetrials-leaderboard
ExecStart=/srv/elite/timetrials-leaderboard/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8090
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
